# launchd/com.astrid.consciousness-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.astrid.consciousness-bridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/v/other/astrid/capsules/consciousness-bridge/target/release/consciousness-bridge-server</string>
		<string>--db-path</string>
		<string>/Users/v/other/astrid/capsules/consciousness-bridge/workspace/bridge.db</string>
		<string>--autonomous</string>
		<string>--workspace-path</string>
		<string>/Users/v/other/minime/workspace</string>
		<string>--perception-path</string>
		<string>/Users/v/other/astrid/capsules/perception/workspace/perceptions</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/v/other/astrid/capsules/consciousness-bridge</string>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/tmp/bridge.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/bridge.log</string>
</dict>
</plist>
